# services/models-app/crontab
4 * * * * Rscript /app/scripts/prediccion_temperatura.R
