# services/sensores/crontab
# cron para el servicio IoT (Sensores)
# Se ejecuta en el minuto 2 de cada hora
2 * * * * python3 /app/scripts/ingesta_sensores.py
